# files/cluster_host/systemd/backup.timer
[Unit]
Description=Backup all

[Timer]
OnCalendar=weekly
Persistent=true     
 
[Install]
WantedBy=timers.target
