# files/cluster_host/systemd/letsencrypt.timer
[Unit]
Description=Renew letsencrypt certificates

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
